# yt-pt-sync.crontab


# if you want to only run between 08:00 and 20:00, replaces the first * by 8-20
0 * * * * /usr/bin/flock -n /run/lock/yt-pt-sync /home/benjamin/yt-pt-sync/yt-pt-sync.php

